# infra/ec2-idle-shutdown/idle-shutdown.service
[Unit]
Description=Idle shutdown check (EC2 auto-stop)
[Service]
Type=oneshot
ExecStart=/usr/local/bin/idle-shutdown.sh
